# lightningos-light/templates/systemd/lnd.service
[Unit]
Description=Lightning Network Daemon (LND)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
User=lnd
Group=lnd
Type=simple
ExecStart=/usr/local/bin/lnd
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/data/lnd /var/log/lnd

[Install]
WantedBy=multi-user.target
